# templates/default/filebeat.service.erb
[Unit]
Description = filebeat Server
After = syslog.target network.target remote-fs.target

[Service]
User = <%= node['hopslog']['user'] %>
PIDFile = <%= node['filebeat']['pid_file'] %>
ExecStart = <%= node['filebeat']['base_dir'] %>/bin/start-filebeat.sh
ExecStop = <%= node['filebeat']['base_dir'] %>/bin/stop-filebeat.sh
LimitNOFILE = 20000

Type = forking

[Install]
WantedBy = multi-user.target
