# templates/default/filebeat.service.erb
[Unit]
Description = filebeat server
After = syslog.target network.target remote-fs.target <%= @deps %>

[Service]
User = <%= @user %>
PIDFile = <%= @pid %>
ExecStart = <%= @exec_start %>
ExecStop = <%= @exec_stop %>
LimitNOFILE = 20000

Type = forking

[Install]
WantedBy = multi-user.target
