# templates/default/kibana.service.erb
[Unit]
Description = kibana Server
After = syslog.target network.target remote-fs.target <%= @deps %>

[Service]
User = <%= node['hopslog']['user'] %>

PIDFile = <%= node['kibana']['pid_file'] %>
ExecStart = <%= node['kibana']['base_dir'] %>/bin/start-kibana.sh
ExecStop = <%= node['kibana']['base_dir'] %>/bin/stop-kibana.sh
LimitNOFILE = 20000

Type = forking

[Install]
WantedBy = multi-user.target
